# batterywarn.timer
[Unit]
Description=Run batterywarn every 5m

[Timer]
OnActiveSec=0s
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target